# debian/naive_backup.service
[Unit]
Description=Backup

[Service]
Type=oneshot
ExecStart=/usr/bin/naive_backup.sh

[Install]
WantedBy=multi-user.target

