# debian/naive-backup.service
[Unit]
Description=Backup

[Service]
Type=oneshot
ExecStart=/usr/bin/naive-backup

